# scripts/service.freebsd
#!/bin/sh

# PROVIDE: opixelserv
# REQUIRE: DAEMON
# KEYWORD: shutdown

: ${opixelserv_enable:=NO}

. /etc/rc.subr

name="opixelserv"
rcvar="opixelserv_enable"
load_rc_config $name

opixelserv_cmd="/usr/local/bin/${name}"
pidfile="/var/run/${name}.pid"
logfile="/var/log/${name}.log"
command="/usr/sbin/daemon"
daemon_flags="-o ${logfile} -P ${pidfile} -r"

# prepend daemon args
opixelserv_flags="${daemon_flags} ${opixelserv_cmd} ${opixelserv_flags}"

run_rc_command "$1"
